# test/TEST-43-CLEAR-ON-FAILURE-DEPS/testsuite.service
[Unit]
Description=Testsuite service

[Service]
ExecStart=/test-dep-removal.sh
Type=oneshot
